# deploy/systemd/certbot-ip-renew.service
[Unit]
Description=Renew the WEBTSTU IP certificate
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/certbot/bin/certbot renew --quiet --deploy-hook "/usr/bin/systemctl reload nginx"
